# service/property_scraper.service
[Unit]
Description=Property scraper (scheduled via scraper.py --schedule)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/property_scraper
ExecStart=/opt/property_scraper/venv/bin/python -u /opt/property_scraper/scraper.py --schedule
StandardOutput=append:/var/log/property_scraper.log
StandardError=append:/var/log/property_scraper.log
Restart=on-failure
User=raspberrypi

[Install]
WantedBy=multi-user.target
